# www/php-81.docker/crontab
* * * * * /usr/local/bin/php /var/www/php-81.docker/commands/example.php